# .config/systemd/user/podman.service.d/override.conf
[Service]
Environment='DOCKER_HOST=unix://%t/podman/podman.sock' 'DOCKER_BUILDKIT=0'
